# autostart/.config/autostart/telegram.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/Telegram
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Telegram
Comment=
